# qbittorrent/rootfs/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bashio

exec /usr/bin/qbittorrent-nox --profile="/config/"
